# atelier_front.service
[Service]
ExecStart=/usr/local/bin/serve -s /root/spark-atelier-leaderboard/front/build

[Unit]
Wants=network-online.target
